# systemd/notammap-notamextractor.timer
[Unit]
Description=Run notam extraction daily.

[Timer]
Unit=notammap-notamextractor.service
OnCalendar=daily

[Install]
WantedBy=multi-user.target